# slack_rtm.service
[Unit]
Description=slack_rtm

[Service]
ExecStart=/usr/bin/python3 /deploy/app/rtm.py
Restart=always

[Install]
WantedBy=multi-user.target
